# conf/supervisor/haproxy.conf
[program:haproxy]
command = fghack haproxy -f /etc/haproxy/haproxy.cfg
priority = 400
autostart = true
