# conf/supervisor/haproxy.conf
[program:haproxy]
command = haproxy
priority = 400
autostart = true
