# conf/supervisor/haproxy.conf
[program:haproxy]
command = fghack haproxy
priority = 400
autostart = true
